# supervisord.conf
[supervisord]
nodaemon=true                ; (start in foreground if true;default false)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[inet_http_server]
port = 127.0.0.1:9001

[supervisorctl]
serverurl=http://127.0.0.1:9001

[program:sshd]
command = /usr/sbin/sshd -D

[program:asterisk]
command = /usr/sbin/asterisk -C /etc/asterisk/asterisk.conf

